# deploy/openronin-backup-hourly.timer
[Unit]
Description=openronin hourly SQLite backup timer
Documentation=https://github.com/openronin/openronin

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
